# seo-snapshot-capture.service
[Unit]
Description=SEO Snapshot Monthly Capture
Documentation=https://seo.theprofitplatform.com.au
After=network.target

[Service]
Type=oneshot
User=avi
Group=avi
WorkingDirectory=/home/avi/projects/seoanalyst/seo-analyst-agent
Environment="PATH=/home/avi/projects/seoanalyst/seo-analyst-agent/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/avi/projects/seoanalyst/seo-analyst-agent/.env
ExecStart=/home/avi/projects/seoanalyst/seo-analyst-agent/venv/bin/python /home/avi/projects/seoanalyst/seo-analyst-agent/capture_monthly_snapshot.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seo-snapshot

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
